# debian/alternc-metrics-basic.timer
[Unit]
Description=Timer for daily history keeping of metrics

[Timer]
OnCalendar=08:00

[Install]
WantedBy=timers.target
